# etl/extract/luxmeterservice/crontab
* * * * * /usr/local/bin/python /app/luxmeter_service.py >> /var/log/luxmeter_service.log 2>&1
