# etl/extract/luxmeterservice/crontab
* * * * * /usr/bin/python3 /app/luxmeter_service.py > /proc/1/fd/1 2>/proc/1/fd/2